# moneymaker.service
[Unit]
Description=Moneymaker Trading Bot
Documentation=https://github.com/Love335/moneymaker
After=network-online.target time-sync.target
Wants=network-online.target time-sync.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/moneymaker/moneymaker/src
ExecStart=/home/moneymaker/moneymaker/venv/bin/python3 /home/moneymaker/moneymaker/src/main.py
ExecStop=/bin/kill -SIGTERM $MAINPID
TimeoutStopSec=30
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=moneymaker
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target